# deploy/roles/docker-registry/templates/systemd/bullseye-web.service
[Unit]
Description=Bullseye Docker Registry
After=network.target

[Service]
Type=simple
User=<%= node[:user] %>
WorkingDirectory=<%= node[:app_path] %>
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose stop
Restart=always

[Install]
WantedBy=multi-user.target
