# deploy/roles/web/templates/systemd/bullseye-web.service
[Unit]
Description=Bullseye Web Server
After=network.target

[Service]
Type=simple
User=<%= node[:user] %>
WorkingDirectory=<%= node[:app_path] %>
ExecStart=/usr/local/bin/bundle exec puma
Restart=always

[Install]
WantedBy=multi-user.target
