# deploy/roles/web/templates/systemd/bullseye-web-batch.timer
[Unit]
Description=Run bullseye web batch process every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
