# deploy/roles/web/templates/systemd/bullseye-web.service
[Unit]
Description=Bullseye Web Server
After=network.target

[Service]
Type=simple
User=<%= node[:user] %>
WorkingDirectory=<%= node[:app_path] %>
Environment=RAILS_ENV=production
ExecStart=/usr/local/bin/bundle exec puma
Restart=always

[Install]
<% if node[:is_vagrant] then %>
WantedBy=multi-user.target vagrant.mount
<% else %>
WantedBy=multi-user.target
<% end %>
